# linux/debian-package/lib/systemd/system/openpath-agent-update.timer
[Unit]
Description=Timer for OpenPath Linux Agent Updates
After=network-online.target

[Timer]
OnBootSec=45min
OnUnitActiveSec=1d
AccuracySec=1h
RandomizedDelaySec=6h
Persistent=true

[Install]
WantedBy=timers.target
