# linux/debian-package/lib/systemd/system/openpath-firewall-restore.service
[Unit]
Description=OpenPath Early-Boot Firewall Restore (fail-closed)
DefaultDependencies=no
Before=network-pre.target
Wants=network-pre.target
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/openpath-firewall-restore.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
