# linux/debian-package/lib/systemd/system/openpath-sse-listener.service
[Unit]
Description=OpenPath SSE Listener for Instant Rule Updates
After=network-online.target dnsmasq.service openpath-dnsmasq.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/openpath-sse-listener.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
