# linux/debian-package/lib/systemd/system/openpath-agent-update.service
[Unit]
Description=OpenPath Linux Agent Update Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/openpath-agent-update.sh
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
